# deployment/med-qa.service
[Unit]
Description=Medical Q&A Portal
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/med-qa-portal/server
ExecStart=/root/.nvm/versions/node/v24.13.0/bin/node index.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target